# deploy/systemd/mission-control-web.service
[Unit]
Description=Mission Control Web
After=network-online.target mission-control-api.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/dilan/mission-control-dashboard
ExecStart=/home/dilan/mission-control-dashboard/scripts/prod/start-web.sh
Restart=always
RestartSec=3

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/dilan/mission-control-dashboard

[Install]
WantedBy=default.target