# deploy/systemd/mission-control-api.service
[Unit]
Description=Mission Control API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/dilan/Project/mission-control-dashboard
ExecStart=/home/dilan/Project/mission-control-dashboard/scripts/prod/start-api.sh
Restart=always
RestartSec=3

# Hardening (lightweight)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/dilan/Project/mission-control-dashboard /home/dilan/note /home/dilan/.openclaw/mission-control

[Install]
WantedBy=default.target